# elbox_WRGND15/templates/aries/wd/storage/udev.rules/50-format-partition.rules
ACTION=="add", KERNEL=="sd[a-z][1-6]", KERNELS=="2-1", RUN+="/etc/scripts/format_partition.sh %k"
